# overlay/rootfs/etc/s6-init/run-image/service/getty/run
#!/bin/execlineb -P
cd /
exec /sbin/getty -L 115200 ttyS0 vt100
