# overlay/rootfs/etc/s6-init/run-image/service/s6-svscan-log/run
#!/bin/execlineb -P
redirfd -w 2 /dev/console
redirfd -w 1 /dev/console

define LOGDIR /var/log/system
s6-envuidgid catchlog
if { mkdir -p $LOGDIR }
if { chmod -R 0750 $LOGDIR }
if { chown -R catchlog:catchlog $LOGDIR }

redirfd -rnb 0 fifo
s6-applyuidgid -Uz --
exec -c --
s6-log -bp --
	T $LOGDIR
